# laoma/project.conf
# 生产环境 supervisor 配置

[program:project]
command = gunicorn project.wsgi_prod -b 127.0.0.1:8000
directory = /root/project/YongDan/project
autostart = true
autorestart = true